# deploy/sverigeklattraren-publisher.timer
[Unit]
Description=Check Sverigeklättraren main for publishable changes

[Timer]
OnBootSec=15s
OnUnitActiveSec=15s
AccuracySec=2s
RandomizedDelaySec=3s
Persistent=true

[Install]
WantedBy=timers.target
